# packaging/debian/lamco-rdp-server.service
[Unit]
Description=Lamco RDP Server
Documentation=https://lamco.ai
After=graphical-session.target
Wants=xdg-desktop-portal.service

[Service]
Type=simple
ExecStart=/usr/bin/lamco-rdp-server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
